# protected/extensions/crontab/crontabs/cronbackup
0 21 * * 5 /usr/bin/php /var/www/mybackup/protected/yiic.php backup --bkId=6
30 21 * * 5 /usr/bin/php /var/www/mybackup/protected/yiic.php backup --bkId=7
