# protected/extensions/crontab/crontabs/cronbackup
46 22 * * * php /var/www/mybk/protected/yiic.php backup --bkId=2
