# protected/extensions/crontab/crontabs/cronbackup
0 16 * * * php /var/www/mybackup/protected/yiic.php backup --bkId=6
30 21 * * 5 php /var/www/mybackup/protected/yiic.php backup --bkId=7
